# deploy/trumpbot-gui.service
[Unit]
Description=TrumpGPT Streamlit Dashboard
After=trumpbot-api.service
Wants=trumpbot-api.service

[Service]
Type=simple
User=arin
WorkingDirectory=/home/arin/trading-bot
ExecStart=/home/arin/trading-bot/venv/bin/streamlit run src/gui/dashboard.py --server.port 8501 --server.address 0.0.0.0 --server.headless true --server.fileWatcherType none
ExecStop=/usr/bin/fuser -k 8501/tcp
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/arin/trading-bot/.env

# Kill all child processes when stopping
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

# Resource limits — Streamlit needs more than 512M
MemoryMax=1G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trumpbot-gui

[Install]
WantedBy=multi-user.target
