# deploy/trumpbot-api.service
[Unit]
Description=TrumpGPT Trading Bot API + Scheduler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=arin
WorkingDirectory=/home/arin/trading-bot
ExecStart=/home/arin/trading-bot/venv/bin/python run_api.py
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/arin/trading-bot/.env

# Resource limits for Pi 4
MemoryMax=2G
CPUQuota=90%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trumpbot-api

[Install]
WantedBy=multi-user.target
